# etc/rc.local
#!/bin/bash
(echo `dumpkeys |grep -i keymaps`; echo keycode 58 = Control) |loadkeys
